# openresty/script/openvue/openvue.cron
#openvue
2 * * * * sh /opt/openresty/openvue/script/openvue/splitlog.sh  > /dev/null 2>&1
8 3 * * * sh /opt/openresty/openvue/script/openvue/log-delete.sh  > /dev/null 2>&1
